# docker/php/supervisor/messenger-worker.conf
; Configuration Supervisor pour le worker Messenger
; A utiliser en production avec Supervisor installe dans le container PHP.
;
; Installation dans le Dockerfile.prod :
;   RUN apt-get update && apt-get install -y supervisor
;   COPY docker/php/supervisor/messenger-worker.conf /etc/supervisor/conf.d/
;
; Demarrage :
;   CMD ["supervisord", "-c", "/etc/supervisor/supervisord.conf"]
;
; Ou en parallele avec PHP-FPM via un script d'entrypoint.

[program:messenger-consume]
command=php /var/www/api/bin/console messenger:consume async --time-limit=3600 --memory-limit=128M
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
startsecs=0
startretries=10
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/messenger-%(process_num)02d.log
stderr_logfile=/var/log/supervisor/messenger-%(process_num)02d-error.log
